# deploy/lawcrm-backend.service
[Unit]
Description=LawCRM FastAPI Backend
After=network.target

[Service]
User=root
WorkingDirectory=/var/www/lawcrm/backend
EnvironmentFile=/var/www/lawcrm/backend/.env
ExecStart=/var/www/lawcrm/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
